# deploy/chameleon-server.service
[Unit]
Description=Chameleon Protocol UDP transport server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/chameleon-protocol
Environment=CHAMELEON_PSK=research-secret
ExecStart=/opt/chameleon-protocol/chameleon-server --address=:9000 --psk=${CHAMELEON_PSK}
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
